# var/lib/dpkg/info/dbus-1.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dbus.service
/lib/systemd/system/dbus.target.wants
/lib/systemd/system/multi-user.target.wants
/usr
/usr/bin
/usr/bin/dbus-cleanup-sockets
/usr/bin/dbus-daemon
/usr/bin/dbus-launch
/usr/bin/dbus-run-session
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/dbus.service
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/dbus.conf
/usr/libexec
/usr/libexec/dbus-daemon-launch-helper
/var
/var/lib
/var/lib/dbus
/lib/systemd/system/dbus.target.wants/dbus.socket
/lib/systemd/system/multi-user.target.wants/dbus.service
