# var/lib/dpkg/info/libpam-runtime.list
/.
/etc
/etc/default
/etc/default/volatiles
/etc/default/volatiles/99_pam
/etc/environment
/etc/pam.d
/etc/pam.d/common-account
/etc/pam.d/common-auth
/etc/pam.d/common-password
/etc/pam.d/common-session
/etc/pam.d/common-session-noninteractive
/etc/pam.d/other
/etc/security
/etc/security/access.conf
/etc/security/faillock.conf
/etc/security/group.conf
/etc/security/limits.conf
/etc/security/limits.d
/etc/security/namespace.conf
/etc/security/namespace.d
/etc/security/namespace.init
/etc/security/pam_env.conf
/etc/security/time.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pam_namespace.service
/usr
/usr/sbin
/usr/sbin/faillock
/usr/sbin/mkhomedir_helper
/usr/sbin/pam_namespace_helper
/usr/sbin/pam_timestamp_check
/usr/sbin/pwhistory_helper
/usr/sbin/unix_chkpwd
/usr/sbin/unix_update
