# var/lib/dpkg/info/weston-init.list
/.
/etc
/etc/default
/etc/default/weston
/etc/pam.d
/etc/pam.d/weston-autologin
/etc/profile.d
/etc/profile.d/weston.sh
/etc/xdg
/etc/xdg/weston
/etc/xdg/weston/weston.ini
/home
/home/weston
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/weston.service
/lib/systemd/system/weston.socket
/lib/systemd/system-preset
/lib/systemd/system-preset/98-weston-init.preset
/usr
/usr/bin
/usr/bin/weston-start
